# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
	if grep -q 'processor.*: 2' /proc/cpuinfo; then
		mask=4
	elif grep -q 'processor.*: 1' /proc/cpuinfo; then
		mask=2
	else
		return
	fi

	for irq in $(grep -E "mt|ra" /proc/interrupts | cut -d: -f1 | sed 's, *,,'); do
		echo "$mask" > "/proc/irq/$irq/smp_affinity"
		[ $mask = 4 ] && mask=8
	done

	modprobe mt_wifi
	ifconfig rax0 up
	ifconfig ra0 up
	ifconfig rai0 up
	ifconfig apclix0 up
	ifconfig apcli0 up
	ifconfig apclii0 up
	/etc/init.d/apcli.sh start
	modprobe mtkhnat

	for eth in /sys/class/net/eth[0-9]*; do
		echo 6 > $eth/queues/rx-0/rps_cpus
	done

	echo "6" > "/sys/class/net/ra0/queues/rx-0/rps_cpus"
 	echo "6" > "/sys/class/net/rai0/queues/rx-0/rps_cpus"
 	echo "6" > "/sys/class/net/rax0/queues/rx-0/rps_cpus"
 	echo "6" > "/sys/class/net/apcli0/queues/rx-0/rps_cpus"
 	echo "6" > "/sys/class/net/apclix0/queues/rx-0/rps_cpus"
 	echo "6" > "/sys/class/net/apclii0/queues/rx-0/rps_cpus"
	
}
