# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
	if grep -q 'processor.*: 2' /proc/cpuinfo; then
		mask=4
	elif grep -q 'processor.*: 1' /proc/cpuinfo; then
		mask=2
	else
		return
	fi

	for irq in $(grep -E "mt|ra" /proc/interrupts | cut -d: -f1 | sed 's, *,,'); do
		echo "2" > "/proc/irq/$irq/smp_affinity"
		
	done

	for iface in $(ls /sys/class/net); do
   [ -e "/sys/class/net/$iface/queues/rx-0/rps_cpus" ] && echo 4 > /sys/class/net/$iface/queues/rx-0/rps_cpus
	done

	modprobe mt_wifi
	ifconfig rax0 up
	ifconfig ra0 up
	ifconfig rai0 up
	ifconfig apclix0 up
	ifconfig apcli0 up
	ifconfig apclii0 up
	/etc/init.d/apcli.sh start
	modprobe mtkhnat
}
